# printer.service
[Unit]
Description=Native printing to reMarkable
Requires=printer.socket

[Service]
ExecStart=/home/root/printer.arm -restart -debug
